# conf/debian/etc/rc.local
#!/bin/sh

pid=$(ps -C imslu.sh -o pid=)
if [ ${#pid} -eq 0 ]; then

    unset pid
    /etc/imslu/scripts/imslu.sh &
else

    logger -p local7.notice -t imslu-scripts "Please wait until the imslu.sh is finished."
    echo "\nPlease wait until the imslu.sh is finished."
fi

# change hashsize on the fly
# https://github.com/jeffmurphy/NetPass/blob/master/doc/netfilter_conntrack_perf.txt
echo 1048576 > /sys/module/nf_conntrack/parameters/hashsize

exit 0
